# terrarium-irc.service
[Unit]
Description=Terra IRC Bot (terrarium-irc harness)
After=network.target
# Note: Requires terrarium-agent HTTP server running on localhost:8080
# The terrarium-agent server should be started separately (it manages vLLM in Docker)

[Service]
Type=simple
User=USER_TO_REPLACE
Group=USER_TO_REPLACE
WorkingDirectory=/path/to/terrarium-irc
Environment="PATH=/path/to/terrarium-irc/venv/bin"
ExecStart=/path/to/terrarium-irc/venv/bin/python /path/to/terrarium-irc/main.py

# Restart on failure
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=terrarium-irc

[Install]
WantedBy=multi-user.target
